# docker/images/interview/config/supervisor.d/conf.d/consumer.ini
[program:consumer]
command=CompileDaemon -exclude-dir=.git -exclude-dir=.idea -exclude-dir=vendor -exclude-dir=docker -exclude-dir=scripts
        -exclude-dir=build --build="go build -o build/consumer cmd/consumer/main.go" --command="./build/consumer" -verbose -graceful-kill
directory=/build
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
